# scripts/memo-app.service
[Unit]
Description=Memo App (Gunicorn for Django)
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/github/memo-app
EnvironmentFile=/home/ec2-user/github/memo-app/.env
ExecStart=/home/ec2-user/github/memo-app/venv/bin/gunicorn config.wsgi:application \
  --bind 127.0.0.1:8001 \
  --workers 3 \
  --access-logfile - \
  --error-logfile -
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
